# trading-agent.service
[Unit]
Description=Hyperliquid AI Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/hyperliquid-trading-agent
EnvironmentFile=/root/hyperliquid-trading-agent/.env
ExecStart=/usr/bin/python3 -m src.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
